# _docker/cron.d/nftspy
5 * * * * root python3 /root/app/main.py > /var/log/nftspy.log 2>&1
